# scripts/zellij-snapshot.service
[Unit]
Description=Take a Zellij session-state snapshot
Documentation=https://github.com/zellij-org/zellij/issues/3752

# Don't run if zellij isn't installed
ConditionPathExists=%h/.cache/zellij

[Service]
Type=oneshot
ExecStart=%h/.local/bin/zellij-snapshot

# Snapshots are best-effort; one failure shouldn't stop the timer
SuccessExitStatus=0 1

# Be nice to the disk
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=default.target
